# scraper.cron
* * * * * python /Users/sparajuli/Desktop/UCB/fun/Stock\ SVM/scraper.py
